# rmills/bin/com.apple.macbe.setSystemSleep.plist
<?xml version="1.0" encoding="UTF-8"?> 
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd"> 
<plist version="1.0"> 
<dict> 
        <key>Label</key> 
        <string>com.apple.macbe.setSystemSleep</string> 
        <key>ProgramArguments</key> 
        <array> 
                <string>/usr/bin/nice</string>
                <string>/usr/sbin/setSystemSleep</string> 
                <string>-setsleep</string> 
                <string>0</string> 
                <string>/Library/LaunchDaemons/com.apple.macbe.setSystemSleep.plist</string> 
        </array> 
        <key>RunAtLoad</key> 
        <true/>
        <key>ServiceDescription</key>
        <string>systemsetup -setsleep 0</string>
        <key>StandardErrorPath</key>
        <string>/var/log/com.apple.macbe.setSystemSleep.log</string>
        <key>StandardOutPath</key>
        <string>/var/log/com.apple.macbe.setSystemSleep.log</string>
</dict> 
</plist>
